# services.d/icinga2-daemon/run
#!/usr/bin/execlineb -P
/usr/sbin/icinga2 daemon --log-level warning
